# conf.d/xrdp.conf
[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman
autostart=true
autorestart=true

[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
autostart=true
autorestart=true
